# deploy_systemd_units/storwyz-postgres-backup.service
[Unit]
Description=Storwyz PostgreSQL verified backup
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
Environment="STORWYZ_DATABASE=superchat_agent"
Environment="STORWYZ_BACKUP_DIR=/var/backups/storwyz/postgres"
Environment="STORWYZ_BACKUP_RETENTION_DAYS=14"
ExecStart=/opt/superchat-ai-agent/web/scripts/backup_postgres.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/storwyz/postgres

[Install]
WantedBy=multi-user.target
